# ansible/roles/keepass/files/keepass.service
[Unit]
Description=Keepass Sync
After=network.target

[Service]
Type=oneshot
ExecStart=/home/radek/bin/keepass-sync
TimeoutStartSec=1min30s

